# .config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=30s
Unit=mbsync.service

[Install]
WantedBy=timers.target
