# .config/systemd/user/pacman-backup.timer
[Unit]
Description=Periodically backup list of installed packages

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
