# .config/systemd/user/pacman-backup.service
[Unit]
Description=Backup list of installed packages

[Service]
ExecStart=/home/maximbaz/bin/pacman-backup.sh
