# .config/systemd/user/qutebrowser-update-useragent.service
[Unit]
Description=Update qutebrowser useragent

[Service]
ExecStart=%h/.local/bin/qutebrowser-update-useragent
