# .config/systemd/user/qutebrowser-update-useragent.timer
[Unit]
Description=Periodically update qutebrowser useragent

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
