# .config/systemd/user/urlwatch.service
[Unit]
Description=Watch selected URLs for changes

[Service]
ExecStart=/usr/bin/urlwatch
