# .config/systemd/user/wallpaper.service
[Unit]
Description=Update wallpaper

[Service]
ExecStart=/home/maximbaz/bin/update-wallpaper.sh
