# .config/systemd/user/waybar-updates.service
[Unit]
Description=Refresh list of outdated packages in waybar

[Service]
ExecStart=%h/.local/bin/waybar-updates refresh
