# etc/systemd/system/openvpn-client@.service.d/override.conf
[Service]
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/maximbaz/.Xauthority"
PrivateTmp=false
ProtectHome=false
