# etc/profile.d/maven.sh
#!/bin/sh

export MAVEN_OPTS="-Xms1g -Xmx12g -XX:PermSize=1g"
