# etc/profile.d/zz_custom.sh
#!/bin/sh

export PATH="$HOME/bin:$PATH:/usr/share/sway/scripts/"
export LIBVA_DRIVER_NAME=iHD
export MOZ_ENABLE_WAYLAND=1
export QT_QPA_PLATFORM=wayland-egl
export QT_WAYLAND_FORCE_DPI=physical
