# etc/systemd/system/iwd.service.d/90-networkmanager.conf
[Unit]
After=systemd-udevd.service
Before=NetworkManager.service
