# etc/systemd/system/backup-repo@.timer
[Unit]
Description=Periodic backup pacman repository to %I

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
