# etc/systemd/system/paccache.timer
[Unit]
Description=Weekly pacman cache cleanup

[Timer]
OnCalendar=weekly
Persistent=true
AccuracySec=12h

[Install]
WantedBy=timers.target
