# etc/systemd/system/paccache.service
[Unit]
Description=Pacman cache cleanup

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
