# etc/systemd/system/reflector.timer
[Unit]
Description=Daily pacman mirrorlist refresh

[Timer]
OnCalendar=daily
Persistent=true
AccuracySec=12h

[Install]
WantedBy=timers.target
