# etc/systemd/system/reflector.timer
[Unit]
Description=Periodic pacman mirrorlist refresh

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
