# etc/systemd/system/system-dotfiles-sync.timer
[Unit]
Description=Periodic system dotfiles sync

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
