# etc/systemd/system/system-dotfiles-sync.service
[Unit]
Description=Sync system dotfiles

[Service]
Type=oneshot
ExecStart=/usr/local/bin/system-dotfiles-sync
